# internal/cmd/setup/pingmon.service
[Unit]
Description=Pingmon monitoring service

[Service]
ExecStart=/usr/bin/env pingmon
Restart=always

[Install]
WantedBy=multi-user.target
